# systemd/openclaw.service
[Unit]
Description=OpenClaw AI Gateway — {{AGENT_NAME}}
After=network.target xvfb.service
Wants=xvfb.service

[Service]
Type=simple
WorkingDirectory={{INSTALL_DIR}}
Environment=DISPLAY=:99
Environment=PATH={{NODE_BIN_DIR}}:{{BUN_BIN_DIR}}:/usr/local/bin:/usr/bin:/bin
Environment=NIMMIT_HOME={{INSTALL_DIR}}
ExecStart={{OPENCLAW_BIN}} gateway start --foreground
Restart=always
RestartSec=5
RestartPreventExitStatus=SIGKILL

[Install]
WantedBy=default.target
